# staging/wade-staging.service
[Unit]
Description=WADE Staging Daemon (full vs light)
After=network-online.target
Wants=network-online.target
ConditionPathExists=/opt/wade/stage_daemon.py

[Service]
Type=simple
User=autopsy
Group=autopsy
# Uncomment if your Splunk UF runs as 'splunk' and you created the group:
SupplementaryGroups=splunk

EnvironmentFile=-/etc/wade/wade.env
Environment=PYTHONUNBUFFERED=1
WorkingDirectory=/opt/wade
ExecStart=/usr/bin/python3 /opt/wade/stage_daemon.py

Restart=on-failure
RestartSec=5s

# Make files group-readable for UF; keep group on new dirs/files
UMask=002

# Sandboxing — allow access to /home and /var/wade for moves & logs
NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=full
ProtectHome=false
ReadWritePaths=/home/autopsy /var/wade

[Install]
WantedBy=multi-user.target
