# scripts/wade-worker.service

[Unit]
Description=WADE Queue Worker
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=autopsy
Group=autopsy
EnvironmentFile=-/etc/wade/wade.env
EnvironmentFile=-/etc/wade/wade-worker.env
WorkingDirectory=/opt/wade
ExecStart=/usr/bin/python3 /opt/wade/wade_worker.py
Restart=on-failure
RestartSec=5s
NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=full
ProtectHome=false
ReadWritePaths=/home/autopsy /var/wade

[Install]
WantedBy=multi-user.target
